# packages/flow-compute/flow_compute-0.1.11a2.tar.gz/flow_compute-0.1.11a2/src/flow/resources/templates/runtime/sections/workload_resume_service.service.j2
[Unit]
Description=Flow Workload Resume Service
After=network-online.target docker.service
Wants=network-online.target
ConditionPathExists=/var/lib/flow/task-state

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/flow-workload-resume.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flow-workload-resume

[Install]
WantedBy=multi-user.target


