# packages/flow-compute/flow_compute-0.1.10rc1.tar.gz/flow_compute-0.1.10rc1/src/flow/resources/templates/runtime/sections/foundrypf.service.j2
[Unit]
Description=Foundry Port Forwarding
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/foundrypf %i
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=foundrypf

[Install]
WantedBy=multi-user.target


