# packages/flow-compute/flow_compute-0.1.10rc1-py3-none-any.whl/flow/resources/templates/runtime/sections/runtime_limit.timer.j2
[Unit]
Description=Flow Runtime Limit Timer
[Timer]
OnBootSec=[[ timer_seconds ]]s
AccuracySec=1min
Persistent=true
[Install]
WantedBy=timers.target


