# packages/flow-compute/flow_compute-0.1.10.tar.gz/flow_compute-0.1.10/src/flow/resources/templates/runtime/sections/runtime_limit.service.j2
[Unit]
Description=Flow Runtime Limit Enforcement
[Service]
Type=oneshot
ExecStart=/usr/local/bin/flow-runtime-cancel.sh
Restart=no


