# packages/simplyprint-duet3d/simplyprint_duet3d-1.3.15.tar.gz/simplyprint_duet3d-1.3.15/simplyprint-duet3d.service
[Unit]
Description=Meltingplot Duet <-> SimplyPrint.io connector
After=network.target
StartLimitIntervalSec=2min
StartLimitBurst=5

[Service]
Type=simple
Restart=always
StartLimitInterval=0
StartLimitBurst=1440
RestartSec=5s
ExecStart=/usr/local/bin/simplyprint-duet3d start
WorkingDirectory=/usr/local
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target